# etc/system/system/vpn-health.service
[Unit]
Description=VPN Health Check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vpn-health-check.sh
